# etc/crontabs/root
# adapt PATH to your needs
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/tmp/usr/bin:/tmp/usr/sbin

# Her 15 dk arayla FTP sunucusuna autowifi log ve interface Bandwidth değerlerini upload et.
# */15 * * * * /sbin/ftpistatistikyedekle

# set cronjob to every 15-->3 minutes:
*/3 * * * * /sbin/autowifipowersave

# interface bandwidth(vnstat) database rutin guncelle
# */2 * * * * /etc/crontabs/vnstat_cron

# per user bandwidth(wrtbwmon) database rutin guncelle
# */2 * * * * /etc/crontabs/wrtbwmon_cron
